# etc/deluge.conf
start on started deluge
stop on stopping deluge

env uid=root
env gid=root
env umask=027

exec start-stop-daemon -S -c $uid:$gid -k $umask -x /usr/bin/deluge-web
